# rpm/remap_per_app/SOURCES/remap-per-app-1.0.0/systemd/remap-per-app.service
[Unit]
Description=Per App Remapping Daemon
After=graphical.target

[Service]
ExecStart=/usr/local/bin/remap-per-app-daemon
Restart=always

[Install]
WantedBy=default.target

